# roles/nginx-certbot/templates/certbot.timer.j2
[Unit]
Description=Certbot Service for domain %i

[Timer]
OnStartupSec=1m
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
