# roles/nginx-certbot/templates/certbot.service.j2
[Unit]
Description=Certbot Service
After=syslog.target network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certbot.sh
